# research_paper_analysis_v3/config/activation/linux/paper-dispatch.path.tmpl
[Unit]
Description=Watch paper trigger file

[Path]
PathChanged={{TRIGGER_FILE}}
Unit=paper-dispatch.service

[Install]
WantedBy=default.target
