# research_paper_analysis_v3/config/activation/linux/paper-dispatch.service.tmpl
[Unit]
Description=Paper machine dispatch once

[Service]
Type=oneshot
WorkingDirectory={{WORKDIR}}
ExecStart={{PYTHON}} -m flatmachines.dispatch_signals --once --resumer config-store --db-path {{DB_PATH}}
